# crontabs/gaspberrypi.latest
#
# Deletes JPG images from directories marked for purging,
# which are at least 24 hours old.
#
1 * * * * /home/pi/svc/bin/purge.sh
